# templates/etc/profile.d/dart.sh.j2
export PATH=$PATH:{{dart_parent_install_dir}}/{{dart_sdk_subdir}}/bin